# files/DEBIAN/postinst
#!/bin/sh

systemctl enable marvind 
systemctl start  marvind
