# files/DEBIAN/postinst
#!/bin/sh

systemctl enable marvind
systemctl daemon-reload
systemctl restart marvind
