# files/DEBIAN/prerm
#!/bin/sh

systemctl stop    marvind 
systemctl disable marvind
